# udev_scripts/20-opi18b20.rules
SUBSYSTEM=="module", DEVPATH=="/module/opi18b20", RUN+="/lib/udev/opi18b20"
